# conf/crontab
#minute hour    mday    month   wday    who         command
*       *       *       *       *       www-data    php /var/www/html/artisan schedule:run
